# home/pi/.config/autostart/nonstoptv.desktop
[Desktop Entry]
Type=Application
Exec=/boot/nonstoptv.py